# deployment/systemd/aquabill-api.service
[Unit]
Description=AquaBill FastAPI Service
After=network.target postgresql.service

[Service]
Type=notify
User=aquabill
WorkingDirectory=/opt/aquabill
Environment="PATH=/opt/aquabill/venv/bin"
ExecStart=/opt/aquabill/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
